# example-supervisor.conf
[program:penguins-observations-application]
user=www-data
directory=/var/www/PENGUINS_EXAMPLE_DIR
command=bash -c ". venv/bin/activate && honcho start"
stopasgroup=true
killasgroup=true
